# pkg/deb/postinst
#!/bin/sh -ex

systemctl unmask cosmian_vm_agent.service
systemctl enable cosmian_vm_agent.service
systemctl stop cosmian_vm_agent.service

systemctl daemon-reload

#DEBHELPER#
