# ansible/roles/ai_runner/templates/cosmian_runner_ai.service.j2
[Unit]
Description=Cosmian Runner AI
Requires=multi-user.target
After=multi-user.target mount_luks.service cosmian_vm_agent.service

[Service]
Type=simple
User=root
ExecStart=cosmian_ai_runner -p 5001
Restart=on-failure
RestartSec=3s
Environment="CONFIG_PATH=/etc/cosmian_vm/config.json"

[Install]
WantedBy=multi-user.target
